# supervisord/celeryd.conf
; ==================================
;  based on this example: https://github.com/celery/celery/blob/3.1/extra/supervisord/celeryd.conf
; ==================================
; This worker should be reserved for diagnoses that need to happen right away,
; rather than for batch jobs that will take a long time to complete.
; The other workers will also handle priority diagnoses if they are not
; performing a batch job
[program:celery_priority]
command=/home/ubuntu/grits-api/grits_api_env/bin/celery worker -A tasks -Q priority --loglevel=INFO
directory=/home/ubuntu/grits-api
user=ubuntu
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; Use different workers for batch jobs so they don't block user submissions.
[program:celery_process]
command=/home/ubuntu/grits-api/grits_api_env/bin/celery worker -A tasks -Q process,priority --loglevel=INFO
directory=/home/ubuntu/grits-api
user=ubuntu
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
killasgroup=true

[program:celery_diagnose]
command=/home/ubuntu/grits-api/grits_api_env/bin/celery worker -A tasks -Q diagnose,priority --loglevel=INFO
directory=/home/ubuntu/grits-api
user=ubuntu
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
killasgroup=true

; If you enable this set the password
; [program:flower]
; command=flower --port=5555 --broker=mongodb://localhost:27017/tasks --basic_auth=user:password
