# supervisord/dashboardd.conf
[program:dashboard]
command=node /home/ubuntu/diagnostic-dashboard/bundle/main.js
environment=METEOR_DB_NAME="$METEOR_DB_NAME",MONGO_URL="$METEOR_MONGO",PORT="$METEOR_PORT",ROOT_URL="$APACHE_URL",MAIL_URL="smtp://localhost"
directory=/home/ubuntu/diagnostic-dashboard/
user=ubuntu
numprocs=1
autostart=true
autorestart=true
startsecs=10
logfile_maxbytes=150MB
logfile_backups=10
loglevel=info
