# supervisord/girderd.conf
[program:girder]
command=/home/ubuntu/girder/girder_env/bin/python -m girder
directory=/home/ubuntu/girder/
user=ubuntu
autostart=true
autorestart=unexpected
