# supervisord/girderd.conf
[program:girder]
command=/home/ubuntu/girder_env/bin/python -m girder
directory=/home/ubuntu/girder/
autostart=true
autorestart=unexpected
