# supervisord/girderd.conf
[program:girder]
command=~/Envs/girder_env/bin/python -m girder
directory=/path/to/girder
autostart=true
autorestart=unexpected
